# test-root.crontab
# start testdns 
*/10 * * * * service named status >/dev/null && cd $HOME/keyroll && ./test-dns.sh >> test-summary.txt
